# agent-config-template/hooks/launchd/com.agent.hermes.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  헤르메스 데몬 상시 실행 (KeepAlive). 데몬이 죽으면 launchd가 자동 재시작한다.
  ProgramArguments 의 __HERMES_START_COMMAND__ 를 실제 헤르메스 실행 명령으로
  교체한 뒤 load 한다. (예: /usr/local/bin/hermes start ...옵션)
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.agent.hermes</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>__HERMES_START_COMMAND__</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>30</integer>
  <key>StandardOutPath</key>
  <string>/tmp/com.agent.hermes.out</string>
  <key>StandardErrorPath</key>
  <string>/tmp/com.agent.hermes.err</string>
</dict>
</plist>
